# skills/AGI-dream-engine/prototype/systemd/arif-dream.timer
[Unit]
Description=Nightly trigger for arifOS Dream Engine memory consolidation

[Timer]
# 04:00 Asia/Kuala_Lumpur daily (20:00 UTC)
OnCalendar=*-*-* 20:00:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
