# skills/agentic-dream-engine/prototype/systemd/arif-dream-phase2.timer
[Unit]
Description=Weekly trigger for arifOS Dream Engine Phase 2 counterfactual rehearsal

[Timer]
# Sunday 02:00 MYT (Saturday 18:00 UTC)
OnCalendar=Sun *-*-* 02:00:00
Persistent=true
# Small randomized delay to avoid I/O spike
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
