# skills/agentic-dream-engine/prototype/systemd/arif-dream-phase2.service
[Unit]
Description=arifOS Dream Engine — Phase 2 Weekly Counterfactual Rehearsal
After=network.target postgresql.service

[Service]
Type=oneshot
# The existing dream-engine runs as root and accesses /root/.hermes/.
# If you migrate to a dedicated user, update WorkingDirectory and permissions.
User=root
Group=root
WorkingDirectory=/root/.hermes/skills/dream-engine
ExecStart=/usr/bin/python3 /root/.hermes/skills/dream-engine/scripts/stage5_rehearse.py --execute
StandardOutput=journal
StandardError=journal
# F1 Kill Switch: rehearsing 20 rules should never take more than 30 minutes.
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
